# rclone/.config/systemd/user/rclone-onedrive.service
[Unit]
Description=rclone mount OneDrive
After=network-online.target

[Service]
Type=notify
# Clear any stale FUSE mount left behind by a previous ungraceful exit (e.g.
# rclone killed across suspend). The leading `-` ignores the error on a clean
# start (nothing to unmount); `-z` is a lazy unmount so a wedged mountpoint is
# detached even if something is stat'ing it. Without this, the Restart=always
# loop wedges forever on "directory already mounted" and hangs ~.
ExecStartPre=-/bin/fusermount -uz %h/OneDrive
ExecStart=/usr/bin/rclone mount onedrive: %h/OneDrive --vfs-cache-mode full --dir-cache-time 30s
ExecStop=/bin/fusermount -u %h/OneDrive
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
