# systemd/system/systemd-hibernate.service.d/20-restore-freeze-session.conf
# Override nvidia-utils' /usr/lib/systemd/system/systemd-hibernate.service.d/
# 10-nvidia-no-freeze-session.conf, which sets SYSTEMD_SLEEP_FREEZE_USER_SESSIONS
# to false. That flag exists for the proprietary nvidia driver's procfs suspend
# interface, which doesn't exist on nvidia-open. With user-session freezing
# disabled, user processes can be in FUSE D-state during the kernel freezer
# pass (statx into sshfs/rclone mounts) and refuse to freeze, causing hibernate
# to abort with "Failed to put system to sleep. Device or resource busy".
# Re-enabling cgroup freezing of user sessions pauses those processes cleanly
# before the kernel freezer runs, avoiding the FUSE deadlock.
# See BOOT_ISSUES.md 2026-04-16.
[Service]
Environment="SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=true"
