# spotify-player/.config/systemd/user/spotify-player-watchdog.timer
[Unit]
Description=Spotify Player Daemon Watchdog Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
