# systemd/system/systemd-suspend.service.d/20-restore-freeze-session.conf
# See sibling drop-in under systemd-hibernate.service.d/ for rationale.
[Service]
Environment="SYSTEMD_SLEEP_FREEZE_USER_SESSIONS=true"
