# hypr/.config/systemd/user/hypr-session-autosave.service
[Unit]
Description=Auto-save Hyprland session (saves periodically + on shutdown)
PartOf=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Save on service stop (triggered by reboot/shutdown/logout)
ExecStop=/usr/bin/hyprsession save

[Install]
WantedBy=graphical-session.target
