# hypr/.config/systemd/user/wallpaper-daily.timer
[Unit]
Description=Fetch daily wallpaper

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
