# hypr/.config/systemd/user/wallpaper-daily.service
[Unit]
Description=Fetch and set daily wallpaper
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hypr-wallpaper fetch
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
